# configs/vrrp.service
# Automatically generated by systemd-sysv-generator

[Unit]
Description=vrrp
After=network.target network-online.target

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
SysVStartPriority=20
ExecStart=/etc/vrrpd/vrrp_on.sh
ExecStop=/usr/bin/killall vrrpd

[Install]
WantedBy=multi-user.target
Alias=vrrp.service
